# sample.timer
[Unit]
Description=Run floodgit daily

[Timer]
OnCalendar=*-*-* 04:00:59
Persistent=true

[Install]
WantedBy=timers.target
